# init/GIH-start_all.conf
description "Green In House start all at init"
start on startup
task
exec /GreenInHouse/script/GIH-start_all.sh